# iosbackupmachine.service
[Unit]
Description=iOS Backup Machine for UDID %I
After=network-online.target usbmuxd.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
# venv python + your main script
ExecStart=/root/iosbackupmachine/bin/python3 /root/iosbackupmachine.py
# env for display and config
Environment=EPD_GPIO_CHIP=/dev/gpiochip3
Environment=EPD_PIN_DC=17
Environment=EPD_PIN_RST=1
Environment=EPD_PIN_BUSY=10
Environment=EPD_SPI_DEV=/dev/spidev3.0
Environment=EPD_SPI_HZ=2000000
Environment=IOSBACKUP_CONFIG=/root/config.yaml
# logs
ExecStartPre=/bin/mkdir -p /var/log/iosbackupmachine
StandardOutput=append:/var/log/iosbackupmachine/autostart.log
StandardError=append:/var/log/iosbackupmachine/autostart.log
# behavior
Restart=no
KillMode=process
TimeoutStopSec=5s

[Install]
WantedBy=multi-user.target
UNIT
