# services/ntp-sync.service
[Unit]
Description=NTP time sync for iOS Backup Machine
After=network-online.target pisugar-server.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Environment=IOSBACKUP_CONFIG=/root/iosbackupmachine/config.yaml
ExecStart=/root/iosbackupmachine/bin/python3 /root/iosbackupmachine/ntp-sync.py
StandardOutput=append:/var/log/iosbackupmachine/ntp-sync.log
StandardError=append:/var/log/iosbackupmachine/ntp-sync.log

[Install]
WantedBy=multi-user.target
