# services/shutdown-display.service
[Unit]
Description=Show owner info on e-ink at shutdown
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecStop=/root/iosbackupmachine/display-shutdown.sh
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
