# services/webui.service
[Unit]
Description=iOS Backup Machine Web UI
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/iosbackupmachine
Environment=IOSBACKUP_CONFIG=/root/iosbackupmachine/config.yaml
ExecStart=/root/iosbackupmachine/bin/python3 /root/iosbackupmachine/webui.py
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
